# cron/crontab
0 */2 * * * /usr/bin/curl http://keeper:3000/parse_job